# firewall.service
[Unit]
Description=Firewall

[Service]
Type=oneshot
ExecStart=/etc/router/firewall

[Install]
WantedBy=multi-user.target